# ops/i5/mamy-proxy.service
[Unit]
Description=MamY Ollama-format proxy to Groq Cloud (with local Ollama fallback)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/mamy-proxy.env
Environment=HOME=/home/optimus
Environment=PATH=/home/optimus/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/usr/bin/python3 /opt/mamy-proxy/proxy.py
Restart=on-failure
RestartSec=5s
User=optimus
WorkingDirectory=/opt/mamy-proxy

[Install]
WantedBy=multi-user.target
